# Docker/php/etc/cron.d/crontab
* * * * *  cd /var/www/html && php artisan schedule:run >> /var/www/html/storage/logs/ubuntu_cron.log
